# scheduled-reboot.timer
[Unit]
Description=Hourly Reboot

[Timer]
OnCalendar=*-*-* *:00:00

[Install]
WantedBy=timers.target
