# scheduled-reboot.service
[Unit]
Description=Hourly Reboot

[Service]
Type=simple
ExecStart=/usr/bin/systemctl --force reboot
