# data/etc/rc.d/rc.local
#!/bin/bash

cp -rp /opt/home/* /home
mkdir /var/log/journal
systemctl restart systemd-journald
