# data/opt/home/display/crontab
DISPLAY=:0

0 22 * * * /opt/vc/bin/tvservice -o
30 7 * * * sudo reboot
0 * * * * killall midori
